# sentinel.service
[Unit]
Description=nexcore-sentinel — Rust SSH brute-force protection
Documentation=https://github.com/nexvigilant/nexcore
After=network.target
Before=sshd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/nexcore-sentinel run -c /etc/nexcore-sentinel/sentinel.toml
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/nexcore-sentinel
ReadOnlyPaths=/var/log/auth.log /etc/nexcore-sentinel
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes

# Needs CAP_NET_ADMIN for iptables, CAP_NET_RAW for raw socket check
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexcore-sentinel

[Install]
WantedBy=multi-user.target
